# threadwatch.timer
# threadwatch.timer — fire every ~5h, ONLY while logged in (no linger), catch up once on login.
# Run standalone: symlink into ~/.config/systemd/user/ then `systemctl --user enable --now threadwatch.timer`
[Unit]
Description=threadwatch 5-hourly knowledge-thread scan (runs only while logged in)

[Timer]
# Every 5 hours at :15, so each scan lands one invocation per Claude 5h usage window
# AND after the hourly claudesync.timer (:00) has refreshed the archive.
OnCalendar=*-*-* 00/5:15:00
# Fire ONCE on next login if a run was missed (no stampede). systemd does not
# allow inline comments after a value, so these notes live on their own lines.
Persistent=true
# Jitter (seconds): less predictable, avoids on-the-dot hammering.
RandomizedDelaySec=300
AccuracySec=1min
Unit=threadwatch.service

[Install]
WantedBy=timers.target
