# threadwatch.timer
# threadwatch.timer — fire ~hourly, ONLY while logged in (no linger), catch up once on login.
# Run standalone: symlink into ~/.config/systemd/user/ then `systemctl --user enable --now threadwatch.timer`
[Unit]
Description=threadwatch hourly knowledge-thread scan (runs only while logged in)

[Timer]
OnCalendar=hourly
# Fire ONCE on next login if a run was missed (no stampede). systemd does not
# allow inline comments after a value, so these notes live on their own lines.
Persistent=true
# Jitter (seconds): less predictable, avoids on-the-dot hammering.
RandomizedDelaySec=300
AccuracySec=1min
Unit=threadwatch.service

[Install]
WantedBy=timers.target
