# threadwatch.service
# threadwatch.service — read-only everywhere except ONE writable dir (~/threadwatch).
# WHY user-level + no linger: runs as you, sees your `claude` OAuth session, and
# "only while logged in" comes free from NOT enabling linger.
# Run standalone: `systemctl --user start threadwatch.service` then
#   `journalctl --user -u threadwatch -o cat -e`
[Unit]
Description=threadwatch read-only knowledge-thread scan

[Service]
Type=oneshot
ExecStart=%h/.local/bin/threadwatch run
TimeoutStartSec=600
# Boot-readiness gate: after a power-loss reboot the Persistent catch-up run fires
# before the network/claude session is ready (observed: started 18:38:40, claude_failed
# 22s later). Block up to 120s for Anthropic reachability, then proceed (logged + retried
# next 5h window). curl works inside this sandbox -- claude makes the same AF_INET calls.
ExecStartPre=/usr/bin/bash -c 'for i in $(seq 1 60); do curl -s -o /dev/null --max-time 3 https://api.anthropic.com/ && exit 0; sleep 2; done; echo "threadwatch: network not ready after 120s, proceeding"; exit 0'

# PATH so the shebang finds `uv` and the script finds `claude` (both in ~/.local/bin).
# UV_NO_CACHE keeps uv from writing a cache under the read-only home.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=UV_NO_CACHE=1
# Guarantee subscription auth ($0/call): claude -p prefers ANTHROPIC_API_KEY over the
# logged-in OAuth session, so strip any inherited key to never make this a paid call.
UnsetEnvironment=ANTHROPIC_API_KEY ANTHROPIC_AUTH_TOKEN

# --- Filesystem: read-only everywhere, ONE writable hole ---
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=%h/disability-email-archive
ReadOnlyPaths=%h/claude-ai-export
ReadOnlyPaths=%h/.claude
ReadOnlyPaths=%h/.claude.json
ReadWritePaths=%h/threadwatch
PrivateTmp=true

# --- Privilege / kernel ---
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true
# ProtectKernelModules / ProtectKernelLogs / ProtectClock are omitted: in a --user
# manager they try to drop CAP_SYS_MODULE/CAP_SYSLOG/CAP_SYS_TIME from the bounding
# set, which needs CAP_SETPCAP and fails the unit (218/CAPABILITIES). The process is
# unprivileged and holds none of those capabilities anyway, so nothing is lost.
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RemoveIPC=true
# CapabilityBoundingSet/AmbientCapabilities are intentionally omitted: a --user
# service is already unprivileged and cannot drop the bounding set without
# CAP_SETPCAP, which fails the unit with 218/CAPABILITIES. NoNewPrivileges covers it.

# --- Syscalls / network (KEEP AF_INET* — the API needs it; DNS is resolved stub on loopback) ---
SystemCallArchitectures=native
# Allowlist @system-service, then deny only @privileged. Do NOT also deny @resources:
# claude uses sched/rlimit syscalls in that group, and a denied call is SIGSYS-killed
# (rc -31 / 31=SIGSYS), not EPERM. (A '~' line is a deny-list.)
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# DO NOT set PrivateNetwork=yes (kills the API). DO NOT set IPAddressDeny=any here
# (cgroup-eBPF is unreliable under --user and would silently break loopback DNS).
# DO NOT set MemoryDenyWriteExecute=true unless a node/V8 claude path is reintroduced.

# --- Resource caps ---
MemoryMax=2G
CPUQuota=80%
TasksMax=256

# --- Logging: wrapper emits one JSON line to stderr -> journal ---
StandardOutput=journal
StandardError=journal
SyslogIdentifier=threadwatch
